# gpib_config.service
[Unit]
Description=Configure GPIB interface

[Service]
Type=oneshot
ExecStart=/usr/bin/gpib_config

[Install]
WantedBy=multi-user.target
